# alpine-consul/rootfs/etc/services.d/go-dnsmasq/run
#!/usr/bin/with-contenv sh

/usr/local/bin/go-dnsmasq --default-resolver --stubzones consul/127.0.0.1:8600
